# debian/cron.daily
#!/bin/sh
test -x /usr/sbin/automysqlbackup && /usr/sbin/automysqlbackup
